# myownservice.service
[Service]
ExecStart=/home/savera/Downloads/2_myownservice.sh

